# .crontab
*/5 * * * * bash /home/germaniii/.config/low-batt-notify.sh
