# etc/systemd/system/x11vnc.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/x11vnc -many -display :0 -no6 -rfbport 5901 -auth /var/run/lightdm/root/:0

[Install]
WantedBy=graphical.target
